# crontab.ini
[program:crontab]
command=/bin/ash /startup
directory=/logs
autostart=true
autorestart=true
startretries=3
stderr_logfile=/logs/supervisor/crontab.err.log
stdout_logfile=/logs/supervisor/crontab.out.log