# nextjs.ini
[program:nextjs]
user=cmnadpl
directory=/apps
command=yarn start
autostart=true
autorestart=true
startretries=3
stderr_logfile=/logs/supervisor/nextjs.err.log
stdout_logfile=/logs/supervisor/nextjs.out.log
environment=HOME="/home/cmnadpl",USER="cmnadpl"
